# deploy/argus.service
[Unit]
Description=Argus Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/argus
EnvironmentFile=/opt/argus/.env
ExecStart=/opt/argus/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
