# deploy/trenchflow.service
[Unit]
Description=Trenchflow Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trenchflow
EnvironmentFile=/opt/trenchflow/.env
ExecStart=/opt/trenchflow/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
